# deploy/playbooks/docker/templates/lib/systemd/system/docker.service.d/exec-start.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon --host={{ docker_host }} --icc=false